# crontab
# crontab

# Set environment variables for the cron jobs
SHELL=/bin/bash
PATH=/home/appuser/.local/bin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

# Log cron jobs to stdout/stderr to be captured by the Docker logger.
# The jobs are staggered to avoid resource contention.

# Full run every 4 hours
30 */4 * * * /home/appuser/.local/bin/python -m src.main full-run >> /proc/1/fd/1 2>/proc/1/fd/2

# Send the digest every hour from 7am to 11am
0 7-11 * * * /home/appuser/.local/bin/python -m src.main send-digest >> /proc/1/fd/1 2>/proc/1/fd/2

# An empty line is required at the end of the file for cron to work correctly.
